# backend/rit-trading.service
[Unit]
Description=RIT Trading Backend Server
After=network.target
Documentation=https://github.com/greenturtle537/rit-trading

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/glitchtech/rit-trading/backend
ExecStart=/usr/bin/perl /glitchtech/rit-trading/backend/server.pl
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rit-trading

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/glitchtech/rit-trading/backend

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PERL5LIB=/glitchtech/rit-trading/backend/local/lib/perl5"

[Install]
WantedBy=multi-user.target
